# service/mesos-exporter/run
#!/bin/sh

exec chpst -U nobody:nogroup /bin/mesos-exporter -master http://leader.mesos:5050/state.json
